# crontab
*/1 * * * * python /app/manage.py tcuperiodicupdate
*/1 * * * * python /app/manage.py tcutimeout
*/1 * * * * python /app/manage.py flushexpiredtokens

