# crontab
* * * * * python /app/manage.py tcuperiodicupdate
* * * * * python /app/manage.py flushexpiredtokens

